# dist/templates/supervisord/app.conf
; ============================================
;  uWSGI with memory limits example for Django
; ============================================


[program:{{PROGRAMNAME}}]
environment=FLAVOR={{FLAVOR}},DJANGO_SETTINGS_MODULE="{{PKGNAME}}.settings_{{FLAVOR}}"
command={{APPDIR}}/.ve/bin/uwsgi
    --socket {{USERDIR}}/run/{{SERVER_NAME}}.sock
    --chmod-socket
    --wsgi-file {{CONFIGABSOLUTENAME}}.wsgi
    --procname-prefix-spaced {{PKGNAME}}-{{FLAVOR}}
    --auto-procname
    --master
    --processes 10
    --no-orphans
    --vacuum
    --log-5xx
    --single-interpreter
    --log-zero
    --log-slow 1000
    --virtualenv {{APPDIR}}/.ve
    --limit-as 300
    --reload-on-as 250
    --reload-on-rss 150
    --forkbomb-delay 0
    --logdate
directory={{APPDIR}}
user={{USERNAME}}
numprocs=1
stdout_logfile={{USERDIR}}/logs/{{PROGRAMNAME}}.log
autostart=true
autorestart=true
startsecs=10
startretries=5
redirect_stderr=true
stopsignal=QUIT
; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs=60
