# dist/templates/supervisord/celerycam.conf
; =======================================
;  celeryd supervisor example for Django
; =======================================

[program:{{PROGRAMNAME}}]
environment=FLAVOR={{FLAVOR}},DJANGO_SETTINGS_MODULE="{{PKGNAME}}.settings_{{FLAVOR}}"
command={{APPDIR}}/.ve/bin/python {{APPDIR}}/src/manage.py celerycam
directory={{APPDIR}}
user={{USERNAME}}
numprocs=1
stdout_logfile={{USERDIR}}/logs/{{PROGRAMNAME}}.log
autostart=true
autorestart=true
startsecs=10
startretries=100
redirect_stderr=true

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs=60

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=996
